# tty/83-ppkb.rules
KERNEL=="pinephone_keyboard", SUBSYSTEM=="module", RUN+="/usr/bin/ppkb_detector.sh"